# config/supervisor.conf
[supervisord]
[program:uwsgi]
command = /usr/bin/uwsgi --ini /etc/app/uwsgi.ini

